# launchd/com.yoshiharu.pomodoro-loop.stop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yoshiharu.pomodoro-loop.stop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/yoshiharu/workspace/YoshiharuTakenaka/pomodoro-loop/stop.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/pomodoro-schedule.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/pomodoro-schedule.err</string>
</dict>
</plist>
