# launchd/local.pomodoro-loop.start.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.pomodoro-loop.start</string>
  <key>ProgramArguments</key>
  <array>
    <string>__SCRIPT_DIR__/start.sh</string>
    <string>--scheduled</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>AbandonProcessGroup</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/pomodoro-schedule.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/pomodoro-schedule.err</string>
</dict>
</plist>
